# ansible/roles/ansible-tvault-contego-extension/ansible-contego-extension-install/files/services/tvault-object-store.service
[Unit]
Description=Tvault Object Store
After=tvault-contego.service
[Service]
User=
Group=
Type=simple
ExecStart=
TimeoutStopSec=20
KillMode=process
Restart=on-failure
[Install]
WantedBy=multi-user.target
