# juju-charms/charm-trilio-data-mover-api/src/files/trilio/trilio-data-mover-api.service
[Unit]
Description=Datamover API service

[Service]
User = root
Group = root
Type = Simple
ExecStart=/usr/bin/python /usr/bin/dmapi-api
KillMode=process
Restart=on-failure

[Install]
WantedBy=multi-user.target
